# hikcentral-middleware.service
[Unit]
Description=HikCentral Middleware Service
Documentation=https://github.com/pokerist/lyve-simple
After=network.target

[Service]
Type=simple
User=nodejs
Group=nodejs
WorkingDirectory=/opt/hikcentral-middleware
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node app.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hikcentral-middleware
KillMode=mixed
TimeoutStopSec=5
KillSignal=SIGTERM

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/hikcentral-middleware

[Install]
WantedBy=multi-user.target
